# rules/susp_registry_persistence.yml
title: Suspicious Registry Run Key Modification
id: f7e8d9c0-b1a2-3456-7890-abcdef123456
status: experimental
description: Detects modification of Windows Registry Run keys, which is a common technique for establishing persistence on a compromised system.
author: CyberBrother
date: 2026-03-19
logsource:
    category: registry_set
    product: windows
detection:
    selection:
        TargetObject|contains:
            - '\SOFTWARE\Microsoft\Windows\CurrentVersion\Run\'
            - '\SOFTWARE\Microsoft\Windows\CurrentVersion\RunOnce\'
            - '\SOFTWARE\WOW6432Node\Microsoft\Windows\CurrentVersion\Run\'
            - '\SOFTWARE\WOW6432Node\Microsoft\Windows\CurrentVersion\RunOnce\'
    filter_trusted:
        Image|endswith:
            - '\msiexec.exe'
            - '\explorer.exe'
    condition: selection and not filter_trusted
falsepositives:
    - Legitimate software installations
    - System updates
level: medium
tags:
    - attack.persistence
    - attack.t1547.001
